# workspace/crontab/laradock
# Do not run crontab on app node (if separate app node and worker node), run only on worker node
# Run "touch -c /etc/crontab /etc/cron.*/*" after starup for preventing hardlink error
* * * * * root php /var/www/mailautoreply-ai-basic/artisan schedule:run >> /var/log/cron.log 2>&1

