# workspace/crontab/laradock
* * * * * laradock php /var/www/mailautoreply-ai-basic/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock for i in {0..10}; do php /var/www/mailautoreply-ai-basic/artisan chatbot:fetch >> /dev/null 2>&1 && sleep 5; done;

#cron job for getting results from tegaki
* * * * * laradock php /var/www/image-to-text-convert/artisan schedule:run >> /dev/null 2>&1

